# llama.cpp/systemd/llama-log-cleanup.service
[Unit]
Description=Clean old llama-server logs from memory

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'find /dev/shm -name "llama-server-*.log" -mmin +10 -delete || true'
